# openwrt/nodogsplash/files/etc/crontabs/root
*/5 * * * * /etc/config/scripts/location_c.sh
59 23 * * * /usr/bin/reboot

