# openwrt/nodogsplash/files/etc/crontabs/root
*/5 * * * * /etc/config/scripts/location_c.sh

